# rsync-home/.config/autostart/make-start.desktop
[Desktop Entry]
Type=Application
Name=autostart by Makefile
Comment=
NoDisplay=false
Terminal=false
Exec=make -C ~/src/mr-realtime start
