# rsync-home/.config/autostart/make-start.desktop
[Desktop Entry]
Type=Application
Name=autostart by Makefile
Comment=
NoDisplay=false
Exec=make -C ~/src start
